# packaging/debian/infmon-frontend.service
[Unit]
Description=InFMon flow-telemetry frontend
Documentation=https://github.com/r12f/InFMon
After=vpp.service network-online.target
Wants=network-online.target
# We cannot Require= / BindsTo= / PartOf= vpp.service because some
# operators run a custom VPP unit name (e.g. vpp-bf3.service) — any
# of those would be a silent no-op there. The frontend instead relies
# on its own retry loop on the binary-API socket (see spec 004 §7.1)
# to handle VPP restarts. Operators who run the canonical vpp.service
# unit and want auto-restart on VPP restart can drop in
# /etc/systemd/system/infmon-frontend.service.d/partof.conf with
#   [Unit]
#   PartOf=vpp.service
# This drop-in path is documented in README.Debian.

[Service]
Type=notify
ExecStart=/usr/bin/infmon-frontend --config /etc/infmon/config.yaml
Restart=on-failure
RestartSec=2s
User=infmon
Group=infmon
# The frontend reads the VPP stats segment and talks to the binary-API
# Unix socket. Both live under /run/vpp by default.
SupplementaryGroups=vpp
RuntimeDirectory=infmon
RuntimeDirectoryMode=0750
StateDirectory=infmon
StateDirectoryMode=0750
LogsDirectory=infmon
LogsDirectoryMode=0750
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
# Stats segment is a hugepages-backed mmap. ReadWritePaths= confines
# access to /dev/hugepages without risking mmap() failures on kernels
# that reject read-only bind mounts for hugepage-backed files.
# Operators who have validated ReadOnlyPaths=/dev/hugepages on their
# target kernel can switch to it via a systemd drop-in for tighter
# hardening.
ReadWritePaths=/dev/hugepages
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
