# deploy/infmon.service
[Unit]
Description=InFMon flow-telemetry frontend
Documentation=https://github.com/r12f/InFMon

# This unit intentionally does NOT use PartOf=vpp.service or BindsTo=vpp.service.
# VPP may run under a custom unit name, and infmon-frontend has its own retry loop
# on the binary-API socket — tying lifecycle to vpp.service would cause unnecessary
# restarts. If you want co-lifecycle with VPP, add PartOf=vpp.service via a drop-in:
#   systemctl edit infmon
#   [Unit]
#   PartOf=vpp.service
After=vpp.service network-online.target
Wants=network-online.target

[Service]
# Type=notify requires the binary to call sd_notify(READY=1). infmon-frontend
# implements sd-notify; if startup takes longer than expected, override with:
#   systemctl edit infmon → TimeoutStartSec=30s
Type=notify

# Default path assumes Debian/package install. For cargo or manual builds the
# binary may be in /usr/local/bin or elsewhere — override via:
#   systemctl edit infmon → ExecStart=/your/path/infmon-frontend ...
ExecStart=/usr/bin/infmon-frontend --config /etc/infmon/config.yaml
Restart=on-failure
RestartSec=2s
User=infmon
Group=infmon

# vpp group is needed for access to the binary-API socket and stats segment
# under /run/vpp.
SupplementaryGroups=vpp
RuntimeDirectory=infmon
RuntimeDirectoryMode=0750
StateDirectory=infmon
StateDirectoryMode=0750
LogsDirectory=infmon
LogsDirectoryMode=0750
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

# Stats segment is hugepages-backed mmap — this grants read access to
# /dev/hugepages. Requires kernel hugepage support; if mmap fails, check
# that hugepages are allocated (see /proc/sys/vm/nr_hugepages).
ReadOnlyPaths=/dev/hugepages

# Capabilities are locked down by default. If you need to bind to a
# privileged port (<1024) or encounter stats-segment permission errors,
# add capabilities back via a drop-in:
#   systemctl edit infmon
#   [Service]
#   CapabilityBoundingSet=CAP_NET_BIND_SERVICE
#   AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
